# private_dot_config/containers/systemd/zeroclaw.container
[Unit]
Description=ZeroClaw Personal AI Assistant
After=network-online.target bifrost.service stocker-mcp.service
Wants=network-online.target

[Container]
Image=ghcr.io/zeroclaw-labs/zeroclaw:v0.8.1-debian
ContainerName=zeroclaw
UserNS=keep-id:uid=65534,gid=65534
PublishPort=42617:42617

# Persistent state + mounted config (config is read-write so onboarding/pairing can persist)
Volume=/storage/zeroclaw/data:/zeroclaw-data:Z

# Resolve host services (bifrost :8080, stoker-mcp :5505) from inside the container
AddHost=host.containers.internal:host-gateway

Environment=ZEROCLAW_ALLOW_PUBLIC_BIND=1
# zc-secrets-injected
Secret=groq_api_key,type=env,target=GROQ_API_KEY
Secret=vk_zeroclaw,type=env,target=ZEROCLAW_API_KEY
Secret=telegram_token_plut01,type=env,target=ZEROCLAW_channels__telegram__default__bot_token


HealthCmd=curl -sf http://localhost:42617/health || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=20s

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=default.target
